# deploy/schedulers/cron/jpad-ops.cron
SHELL=/bin/bash
PATH=/usr/local/bin:/usr/bin:/bin

JPAD_ROOT=/srv/jpad
JPAD_ENV_FILE=/srv/jpad/.env
JPAD_LOCK_DIR=/var/lock/jpad
SEARCH_INDEX_JOB_LIMIT=50
AUDIT_LOG_DELIVERY_LIMIT=50
ATTACHMENT_SECURITY_RESCAN_LIMIT=50

# /etc/cron.d/jpad-ops 템플릿
# 적용 전:
# 1. /var/log/jpad 디렉터리 생성
# 2. deploy/scripts/run-scheduled-job.sh 실행 권한 부여
# 3. 경로/사용자(jpad) 확인

*/5 * * * * jpad $JPAD_ROOT/deploy/scripts/run-scheduled-job.sh semantic-index-jobs >> /var/log/jpad/semantic-index-jobs.log 2>&1
*/2 * * * * jpad $JPAD_ROOT/deploy/scripts/run-scheduled-job.sh audit-log-deliveries >> /var/log/jpad/audit-log-deliveries.log 2>&1
20 * * * * jpad $JPAD_ROOT/deploy/scripts/run-scheduled-job.sh attachment-security-rescan >> /var/log/jpad/attachment-security-rescan.log 2>&1
15 3 * * * jpad $JPAD_ROOT/deploy/scripts/run-scheduled-job.sh retention >> /var/log/jpad/retention.log 2>&1
0 2 * * * jpad $JPAD_ROOT/deploy/scripts/run-scheduled-job.sh backup >> /var/log/jpad/backup.log 2>&1
30 4 * * 0 jpad $JPAD_ROOT/deploy/scripts/run-scheduled-job.sh restore-drill >> /var/log/jpad/restore-drill.log 2>&1

# 전체 재색인은 무거우므로 주 1회 정도로 제한하는 편이 안전합니다.
0 1 * * 0 jpad $JPAD_ROOT/deploy/scripts/run-scheduled-job.sh semantic-reindex >> /var/log/jpad/semantic-reindex.log 2>&1
